# src/dist/omero-ms-pixel-buffer.service
[Unit]
Description=OMERO pixel data microservice server
Requires=multi-user.target
After=multi-user.target
After=network.service

[Service]
Type=simple
Environment="JAVA_OPTS=-Dlogback.configurationFile=/opt/omero/OMERO.ms/omero-ms-pixel-buffer/current/conf/logback.xml"
WorkingDirectory=/opt/omero/OMERO.ms/omero-ms-pixel-buffer/current
ExecStart=/opt/omero/OMERO.ms/omero-ms-pixel-buffer/current/bin/omero-ms-pixel-buffer
User=omero
Group=omero
Restart=no
LimitNPROC=16384
LimitNOFILE=16384

[Install]
WantedBy=multi-user.target
